# templates/powerbt.service
[Unit]
Description=powerbt

[Service]
ExecStart=/powerbt/main.py

[Install]
WantedBy=multi-user.target