# templates/powersw.service
[Unit]
Description=powersw

[Service]
ExecStart=/powersw/main.py

[Install]
WantedBy=multi-user.target